# bin/com.m3memory.cognitiveloop.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  m3-memory Cognitive Loop — macOS launchd agent template.

  Installed by bin/install_schedules.py, which substitutes [M3_MEMORY_ROOT]
  and [M3_PYTHON] and writes the result to
  ~/Library/LaunchAgents/com.m3memory.cognitiveloop.plist.

  NOTE: ProgramArguments deliberately does NOT pass the background flag. launchd
  manages the process directly; the loop must run in the foreground or
  KeepAlive would lose track of the real worker after the double-fork.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.m3memory.cognitiveloop</string>
    <key>ProgramArguments</key>
    <array>
        <string>[M3_PYTHON]</string>
        <string>[M3_MEMORY_ROOT]/bin/m3_cognitive_loop.py</string>
        <string>--interval</string>
        <string>300</string>
        <string>--log-file</string>
        <string>[M3_MEMORY_ROOT]/logs/cognitive_loop.log</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <!-- KeepAlive restarts only on an ABNORMAL (crash/signal) exit, NOT on a
         clean exit. A clean exit 4 = EXIT_ALREADY_RUNNING (another instance
         holds the single-instance lock); unconditional KeepAlive=true would
         respawn-loop on it. Crashed-only is the launchd idiom for "restart if it
         died, leave it alone if it chose to exit". -->
    <key>KeepAlive</key>
    <dict>
        <key>Crashed</key>
        <true/>
    </dict>
    <key>StandardOutPath</key>
    <string>[M3_MEMORY_ROOT]/logs/cognitive_loop.log</string>
    <key>StandardErrorPath</key>
    <string>[M3_MEMORY_ROOT]/logs/cognitive_loop.log</string>
</dict>
</plist>
